# systemd/animaltracker.service.d/memory.conf
# Memory backstops for animaltracker.service
#
# Without these, a slow leak (e.g. the in-memory event-frame buffer that
# caused the 2026-05-02 OOM) accumulates until the kernel OOM-killer
# SIGKILLs the entire process -- which discards in-flight sidecar writes
# and orphan-clip metadata. With these, the cgroup throttles at
# MemoryHigh and the kernel kills *inside the cgroup* at MemoryMax
# (logged to journal), letting systemd Restart= cycle us cleanly.
#
# Sized for a 70 GiB host (74 GB VM after firmware reservation):
#   - MemoryHigh=58G  -> kernel reclaims aggressively above this; pages
#                        the warning into journal
#   - MemoryMax=64G   -> hard ceiling; cgroup OOM-kill above this with
#                        memory.events recorded in journal
#
# Install:
#   sudo cp -r systemd/animaltracker.service.d /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl restart animaltracker

[Service]
MemoryAccounting=yes
MemoryHigh=58G
MemoryMax=64G
# Don't bias us toward being killed by the global OOM killer; let the
# cgroup limits trigger first so the kill is contained and logged.
OOMScoreAdjust=-100
