# systemd/detector@.service
[Unit]
Description=Animal Tracker detector (%i)
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/mnt/wildlife_ssd

[Service]
Type=simple
# Update User and paths to match your system
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/animaltracker
EnvironmentFile=/home/YOUR_USERNAME/animaltracker/config/secrets.env
ExecStart=/home/YOUR_USERNAME/animaltracker/.venv/bin/python -m animaltracker.cli \
    --config /home/YOUR_USERNAME/animaltracker/config/cameras.yml \
    run --camera %i
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
