# systemd/ssd-cleaner.service
[Unit]
Description=Animal Tracker SSD cleanup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# Update User and paths to match your system
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/animaltracker
EnvironmentFile=/home/YOUR_USERNAME/animaltracker/config/secrets.env
ExecStart=/home/YOUR_USERNAME/animaltracker/.venv/bin/python -m animaltracker.cli cleanup \
    --config /home/YOUR_USERNAME/animaltracker/config/cameras.yml
